# initfiles/bash/etc/systemd/system/hpcc-environment-monitor.service.in
[Unit]
Description=regeneration of hpccsystems-platform.target
Conflicts=hpccsystems-platform.target

[Service]
Type=oneshot
ExecStart=@ADMIN_PATH@/generate-hpccsystems-target.sh
