# initfiles/bash/etc/systemd/system/hpcc-init@instance.service.in
[Unit]
Description=hpcc-systems platform component %i

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/hpcc-init -c %i start
ExecStop=/etc/init.d/hpcc-init -c %i stop
WorkingDirectory=@RUNTIME_PATH@
StandardOutput=journal+console

[Install]
DefaultInstance=mydafilesrv
