# initfiles/bash/etc/systemd/system/hpcc-environment-monitor.path.in
[Unit]
Description=Triggers the service that regenerates hpccsystems-platform.target

[Path]
PathModified=@CONFIG_DIR@/@ENV_XML_FILE@

[Install]
WantedBy=multi-user.target
