# initfiles/bash/etc/systemd/system/sparkthor@instance.service.in
[Unit]
Description=Spark cluster for hpcc connectivity
After=network.target remote-fs.target hpcc-conf@%i.service dafilesrv.service
Requires=hpcc-conf@%i.service

[Service]
Type=oneshot
User=@RUNTIME_USER@
Group=@RUNTIME_GROUP@
PermissionsStartOnly=false
WorkingDirectory=@RUNTIME_PATH@/%i
ExecStart=@ADMIN_PATH@/sparkthor.sh start %i
ExecStop=@ADMIN_PATH@/sparkthor.sh stop %i
RemainAfterExit=yes

[Install]
DefaultInstance=myspark
