# initfiles/bash/etc/systemd/system/thorslave@instance.service.in
[Unit]
Description=%i
After=network.target remote-fs.target @dependencies@

[Service]
Type=forking
User=@RUNTIME_USER@
Group=@RUNTIME_GROUP@
PermissionsStartOnly=false
PIDFile=@PID_PATH@/%i.pid
EnvironmentFile=@RUNTIME_PATH@/thorslaves/%i.cfg
ExecStart=@EXEC_PATH@/@processName@ @processArgs@

[Install]
DefaultInstance=my@component@
